# setup/external-web/DEBIAN/postinst
#!/usr/bin/env bash 
set -euxo pipefail

systemctl enable web-server 
systemctl start web-server